# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:dev-server]
command=node /var/www/project/dev-server.js
directory=/var/www/project
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/dev-server.err.log
stdout_logfile=/var/log/supervisor/dev-server.out.log
priority=10

[program:build-watcher]
command=node /var/www/project/build.js --watch
directory=/var/www/project
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/build-watcher.err.log
stdout_logfile=/var/log/supervisor/build-watcher.out.log
priority=20
